# roles/deploy_files/files/573-mobula.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="{{install_dir}}/bin/udev_hook.sh '%E{INTERFACE}'"
